# src/scheduler/deployment/scheduler.service
[Unit]
Description=Trading Platform Scheduler Service
Documentation=https://github.com/your-org/trading-platform
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-platform
ExecStart=/opt/trading-platform/.venv/bin/python -m src.scheduler.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scheduler-service

# Environment variables
Environment=TRADING_ENV=production
Environment=LOG_LEVEL=INFO
Environment=SCHEDULER_MAX_WORKERS=10
Environment=NOTIFICATION_SERVICE_URL=http://localhost:8000
Environment=DATA_CACHE_ENABLED=true

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trading-platform/logs
ReadWritePaths=/opt/trading-platform/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target