# bin/install-ibkr/ibgw-paper-watchdog.service
[Unit]
Description=Restart ib-gateway-paper (PAPER only) if stuck in socat/Xvfb bind-retry loop

[Service]
Type=oneshot
ExecStart=/opt/apps/e-trading/bin/install-ibkr/ibgw_paper_watchdog.sh
