# bin/install-ibkr/ibgw-paper-watchdog.timer
[Unit]
Description=Run ibgw-paper-watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=ibgw-paper-watchdog.service

[Install]
WantedBy=timers.target
